# com.claude.reviewer-prs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.reviewer-prs</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/hibiki.tatsuno/claude-widget/update-reviewer-prs.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/reviewer-prs.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/reviewer-prs.error.log</string>
</dict>
</plist>
