# com.claude.usage-cache.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.usage-cache</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/hibiki.tatsuno/.claude/cache/update-usage-cache.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/hibiki.tatsuno/.claude/cache/update.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hibiki.tatsuno/.claude/cache/update.log</string>
</dict>
</plist>
